# systemd/user/darkman.service
[Unit]
Description=Darkman service

[Service]
ExecStart=/usr/bin/darkman run
Restart=on-failure

[Install]
WantedBy=default.target
