# systemd/user/kanata.service
[Unit]
Description=Kanata keyboard remapper
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/kanata_cmd_allowed --cfg %h/.config/kanata/config.kbd
Restart=on-failure

[Install]
WantedBy=graphical-session.target
